# deploy/mex-self-update.service
[Unit]
Description=MeX Next self-update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=HOME=/root
Environment=MEX_NEXT_DIR=/opt/mex-next
Environment=MEX_NEXT_REF=main
Environment=BOT_SERVICE=mex-bot.service
EnvironmentFile=-/etc/mex/self-update.env
WorkingDirectory=/opt/mex-next
ExecStart=/bin/bash /opt/mex-next/scripts/self-update.sh
StandardOutput=journal
StandardError=journal
SuccessExitStatus=0
