# deploy/timers/mex-proactive-nudge-stale-target.service.template
[Unit]
Description=MeX Next stale target proactive nudge ({ACCOUNT_ID})

[Service]
Type=oneshot
EnvironmentFile=/etc/mex/{ACCOUNT_ID}.env
WorkingDirectory=/opt/mex-next
ExecStart=/usr/bin/doppler run --project xops-{ACCOUNT_ID} --config prd -- /usr/bin/node /opt/mex-next/dist/scripts/cron-proactive-nudge.js --account-id {ACCOUNT_ID} --kind stale_target_review
