# deploy/timers/mex-reactions-poll.service.template
[Unit]
Description=MeX Next inbound reactions poll ({ACCOUNT_ID})
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=HOME=/root
EnvironmentFile=/etc/mex/{ACCOUNT_ID}.env
WorkingDirectory=/opt/mex-next
ExecStart=/usr/bin/doppler run --project xops-{ACCOUNT_ID} --config prd -- /usr/bin/node /opt/mex-next/dist/main.js --account-id {ACCOUNT_ID} --task reactions-poll
StandardOutput=journal
StandardError=journal
